# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/bin/bash -c "/usr/sbin/php-fpm7.0 -F"

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx"
